# esp32-scanner.service
[Unit]
Description=ESP32 Scanner Server
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/esp32-scanner
ExecStart=/usr/bin/node /home/ubuntu/esp32-scanner/server.js
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=esp32-scanner
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
